# config/php7.4/supervisord/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
environment=APACHE_LOG_DIR=/var/log/httpd
redirect_stderr=true
